# systemd/epd-status-dashboard.timer
[Unit]
Description=Render EPD server status dashboard every 3 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=3min
AccuracySec=30s
Persistent=true
Unit=epd-status-dashboard.service

[Install]
WantedBy=timers.target
